# resources/systemd/restart-avahi.timer
[Unit]
Description=timer to restart avahi-daemon.service

[Timer]
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target
